# initscript/Mandriva/networkmanager.in
#!/bin/sh
#
# NetworkManager:   NetworkManager daemon
#
# chkconfig: - 98 02
# description:  This is a daemon for automatically switching network \
#               connections to the best available connection.
#
# processname: NetworkManager
# pidfile: /var/run/NetworkManager/NetworkManager.pid
#
### BEGIN INIT INFO
# Provides: networkmanager
# Required-Start: $network
# Required-Stop: $network
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: Daemon for automatically switching to best network connection.
# Description: This is a daemon for automatically switching network
#              connections to the best available connection.
### END INIT INFO

prefix=@prefix@
exec_prefix=@prefix@
sbindir=@sbindir@

NETWORKMANAGER_BIN=${sbindir}/NetworkManager

# Sanity checks.
[ -x $NETWORKMANAGER_BIN ] || exit 1

# Source function library.
. /etc/rc.d/init.d/functions

# Source network configuration
. /etc/sysconfig/network

# so we can rearrange this easily
processname=NetworkManager
servicename=NetworkManager
pidfile=/var/run/NetworkManager/NetworkManager.pid

RETVAL=0

start()
{
	echo $"Setting network parameters... "
	sysctl -e -p /etc/sysctl.conf >/dev/null 2>&1

	if [ ! -e /var/lock/subsys/named ]; then
		service named start >/dev/null 2>&1
	fi

	echo -n $"Starting NetworkManager daemon: "
	daemon --check $servicename $processname --pid-file=$pidfile
	RETVAL=$?
	echo
	if [ -n "${NETWORKWAIT}" ]; then 
		[ -z "${LINKDELAY}" ] && LINKDELAY=10
		echo -n $"Waiting for network..."
		nm-online -q --timeout=$LINKDELAY || nm-online -q -x --timeout=30
		[ "$?" = "0" ] && success "network startup" || failure "network startup"
		echo
		[ -n "${NETWORKDELAY}" ] && /bin/sleep ${NETWORKDELAY}
	fi
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/$servicename
}

stop()
{
	echo -n $"Stopping NetworkManager daemon: "
	killproc -p $pidfile $servicename
	RETVAL=$?
	echo
	if [ $RETVAL -eq 0 ]; then
		rm -f /var/lock/subsys/$servicename
		rm -f $pidfile
	fi
}

# See how we were called.
case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	status)
		status -p $pidfile $processname
		RETVAL=$?
		;;
	restart)
		stop
		start
		;;
	condrestart)
		if [ -f /var/lock/subsys/$servicename ]; then
			stop
			start
		fi
		;;
	*)
		echo $"Usage: $0 {start|stop|status|restart|condrestart}"
		;;
esac
exit $RETVAL
