# initscript/Gentoo/NetworkManagerDispatcher.in
#!/sbin/runscript
# Copyright (c) 2008 Saleem Abdulrasool <compnerd@compnerd.org>
# Distributed under the terms of the GNU General Purpose License v2
# $Header: $

depend() {
   need NetworkManager
}

start() {
   ebegin "Starting NetworkManagerDispatcher"
      start-stop-daemon --start --quiet --pidfile /var/run/NetworkManagerDispatcher.pid --exec /usr/sbin/NetworkManagerDispatcher
      eend $?
}

stop() {
   local retval

   ebegin "Stopping NetworkManagerDispatcher"
   start-stop-daemon --stop --quiet --pidfile /var/run/NetworkManagerDispatcher.pid
   retval=$?
   eend ${retval}

   [ -e /var/run/NetworkManagerDispatcher.pid ] && rm -f /var/run/NetworkManagerDispatcher.pid
}

# vim: set ft=gentoo-init-d ts=3 sw=3 et:
