# initscript/SUSE/networkmanager.in
#! /bin/sh

### BEGIN INIT INFO
# Provides:       networkmanager
# Required-Start: dbus haldaemon
# Default-Start:  3 5
# Default-Stop:
# Description: NetworkManager, a network and wireless connection manager
### END INIT INFO

sbindir=@sbindir@

NETWORKMANAGER_BIN=${sbindir}/NetworkManager
test -x $NETWORKMANAGER_BIN || exit 5

. /etc/rc.status
rc_reset

case "$1" in
	start)
		checkproc $NETWORKMANAGER_BIN
		if [ $? = 0 ]; then
			echo -n "NetworkManager already running"
			rc_status -v
			rc_exit
		fi
		echo -n "Starting the NetworkManager"
		startproc $NETWORKMANAGER_BIN
		rc_status -v
		;;
	stop)
		echo -n "Shutting down the NetworkManager"
		killproc -TERM $NETWORKMANAGER_BIN
		rc_status -v
		;;
	restart)
		$0 stop
		$0 start
		rc_status
		;;
	status)
		echo -n "Checking for the NetworkManager: "
		checkproc $NETWORKMANAGER_BIN
		rc_status -v
		;;
	*)
		echo "Usage: $0 {start|stop|status|restart}"
		exit 1
		;;
esac

