# data/NetworkManager-wait-online.service.in
[Unit]
Description=Network Manager Wait Online
Requisite=NetworkManager.service
After=NetworkManager.service
Wants=network.target
Before=network.target network-online.target

[Service]
Type=oneshot
ExecStart=@bindir@/nm-online -s -q --timeout=30

[Install]
WantedBy=multi-user.target
