# data/NetworkManager.service.in
[Unit]
Description=Network Manager
Wants=network.target
Before=network.target @DISTRO_NETWORK_SERVICE@

[Service]
Type=dbus
BusName=org.freedesktop.NetworkManager
ExecStart=@sbindir@/NetworkManager --no-daemon
# NM doesn't want systemd to kill its children for it
KillMode=process

[Install]
WantedBy=multi-user.target
Alias=dbus-org.freedesktop.NetworkManager.service
Also=NetworkManager-dispatcher.service

