# data/NetworkManager-dispatcher.service.in
[Unit]
Description=Network Manager Script Dispatcher Service

# Order the dispatcher before NetworkManager. While dispatcher
# is D-Bus activate (and not intended to be explicitly wanted by
# another service/target), the ordering dependency matters during
# shutdown. We want first NetworkManager to be stopped.
Before=NetworkManager.service

[Service]
Type=dbus
BusName=org.freedesktop.nm_dispatcher
ExecStart=@libexecdir@/nm-dispatcher

# We want to allow scripts to spawn long-running daemons, so tell
# systemd to not clean up when nm-dispatcher exits
KillMode=process

[Install]
Alias=dbus-org.freedesktop.nm-dispatcher.service

