# data/NetworkManager-wait-online.service.in
[Unit]
Description=Network Manager Wait Online
Requisite=NetworkManager.service
After=NetworkManager.service
Before=network-online.target

[Service]
Type=oneshot
ExecStart=@bindir@/nm-online -s -q --timeout=30

[Install]
WantedBy=network-online.target
