# clients/cloud-setup/nm-cloud-setup.service.in
[Unit]
Description=Automatically configure NetworkManager in cloud

[Service]
Type=oneshot
ExecStart=@libexecdir@/nm-cloud-setup

#Environment=NM_CLOUD_SETUP_LOG=TRACE
